# saiborg.service
[Unit]
Description=Saiborg Slack Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/saiborg
Environment="PATH=/usr/bin:/usr/local/bin"
ExecStart=/usr/bin/python3 /root/saiborg/app.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

